# systemd/telegram-cleanup.service
[Unit]
Description=Telegram Research Agent — DB Cleanup
After=network.target

[Service]
Type=simple
User=oc_you
Group=oc_you
WorkingDirectory=/srv/openclaw-you/workspace/telegram-research-agent
EnvironmentFile=/srv/openclaw-you/.env
Environment=PYTHONPATH=/srv/openclaw-you/workspace/telegram-research-agent
ExecStart=/srv/openclaw-you/venv/bin/python3 src/main.py cleanup
Restart=no
StandardOutput=journal
StandardError=journal
SyslogIdentifier=telegram-cleanup
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/srv/openclaw-you/workspace/telegram-research-agent/data

[Install]
WantedBy=multi-user.target
